# zoneminder-supervisord.conf
[program:zoneminder]
command = /usr/bin/zmpkg.pl start
pidfile = /var/run/zoneminder/zm.pid
autostart=false
autorestart=true
startsecs = 0
